# minio-proxy/scripts/minio-proxy.ini
[program:minio-proxy]
command=/Users/edtenz/go/bin/minio-proxy -endpoint localhost:9090 -key minioadmin -secret minioadmin -port 10086 -username admin -password admin
process_name=%(program_name)s
directory=/tmp
autostart=true
autorestart=true
startsecs=5
exitcodes=0,2
user=edtenz
redirect_stderr=true
stdout_logfile=/opt/homebrew/var/log/minio-proxy.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5